# systemd/type1-monitor.service
[Unit]
Description=Type1 Monitor FastAPI Backend
After=network.target mysql.service

[Service]
WorkingDirectory=/srv/type1-monitor/backend
ExecStart=/srv/type1-monitor/backend/venv/bin/python3 -m uvicorn main:app --host 127.0.0.1 --port 8002
Restart=always
EnvironmentFile=/srv/type1-monitor/backend/.env

[Install]
WantedBy=default.target
